# examples/app_supervisord.ini
# path: /etc/supervisor.d/app.ini
[program:http_server]
command=python3 -m http.server 8080
autostart=true
autorestart=true
stdout_logfile=/var/log/stdout.log
stderr_logfile=/var/log/stderr.log
